# meta-aws-demos/recipes-core/rauc/files/rauc.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/rauc --mount=/run/rauc/mnt -d service
