# meta-aws-demos/recipes-iot/aws-iot-greengrass/greengrass-lite/ggl.gg_fleetprovisioning.service
[Unit]
After=ggl.core.ggconfigd.service
Requires=ggl.core.ggconfigd.service

[Service]
ExecStart=/usr/bin/fleetprovisioning
